# config/systemd/user/plasma-powerdevil.service.d/override.conf
[Service]
Environment="QT_QPA_PLATFORM=xcb"
Environment="KSCREEN_BACKEND_INPROCESS=1"
Environment="KSCREEN_BACKEND=XRandR"
